# ohd/_lib_systemd_system_ Files/ohd.service
[Unit]
Description=Garage door monitor. ohd.service
After=network.target

[Service]
ExecStart=/home/Your_User/ohd/ohd.py
WorkingDirectory=/home/Your_User/ohd/
StandardOutput=syslog
StandardError=syslog
User=Your_User
ExecStop = /bin/kill -2 $MAINPID

[Install]
WantedBy=multi-user.target
